# template/supervisord.d/sixpack-web.ini
[program:sixpack-web]
command = /home/sixpack/sixpack/bin/sixpack-web
autorestart = true
stdout_logfile = /var/log/supervisor/sixpack-web.log
redirect_stderr = true
environment = PYTHONPATH="."
directory = /home/sixpack/sixpack
